# supervisor-nowgoal.conf
[program:nowgoalcrawler]
command=/opt/app/run 1
numprocs=10
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true